# overlay/system/etc/systemd/system/usb-moded.service.d/usb_moded_wait_usb_init.conf
[Service]
ExecStartPre=/bin/bash /etc/systemd/system/usb-moded.service.d/wait_for_usb_init.sh
