# overlay/system/usr/lib/systemd/system/usb-pre.service
[Unit]
Description=exynosUSB

[Service]
Type=oneshot
ExecStart=/usr/libexec/usb-moded/usb-exy.sh pre
