# overlay/system/usr/lib/systemd/system/usb-post.service
[Unit]
Description=exynosUSB

[Service]
Type=notify
ExecStart=/usr/libexec/usb-moded/usb-exy.sh post
